# deploy/systemd/rwoo-api.service.d/shared-state.conf
[Service]
EnvironmentFile=-/etc/rwoo/rwoo-release.env
ReadWritePaths=/var/lib/rwoo /var/cache/rwoo
